# profile.d/functions.sh
#!/bin/bash
#==========================================================================
#
#         FILE:  functions.sh
#
#  DESCRIPTION:  a collection of functions
#
#==========================================================================

[[ "$-" =~ i ]] || return


