# profile.d/key.sh
GPG_AGENT_INFO=/tmp/gpg-oXgdOm/S.gpg-agent:28802:1; export GPG_AGENT_INFO;
